# docker/cron.d/laravel-cron
# 日本時間の朝8時に統計情報を更新（直接実行）
0 8 * * * root cd /var/www/html && php artisan stats:update >> /var/log/cron.log 2>&1
