# docker/cron.d/laravel-cron
#!/bin/sh
# 毎分 Laravel スケジューラを実行（JST）
# これにより routes/console.php の everyTwoHours() が正しく反映されます
* * * * * root cd /var/www/html && TZ=Asia/Tokyo php artisan schedule:run >> /var/log/cron.log 2>&1 || echo "$(date): schedule:run failed" >> /var/log/cron.log
